# src/etc/uci-defaults/rtorrent
#!/bin/sh
rm -fr /tmp/luci-indexcache /tmp/luci-modulecache
exit 0
